# example/test.cron
* 2 3 4 5 mplayer ~/music/Recorded/jeans-made-out-of-sweat-pants.mp3
