# test.cron
# This file was automatically generated by crontab.go
# Do not manually make changes to this file, they will be overwriten
#
# For more information visit http://www.github.com/knickers/crontab
# Also see the Linux manual pages of crontab(5) and cron(8)

*/2 3 4,7 12 * mplayer ~/music/Recorded/jeans-made-out-of-sweat-pants.mp3
#This is a comment
* 3 4 12 * echo hello world
#echo comment
30 */3 * * */3 ls ~/go
#
